# templates/etc/systemd/system/update-route53-arr.service.j2
[Unit]
Description=Update A RR on Route53
After=network.target

[Service]
Type=simple
RemainAfterExit=yes
ExecStart=/usr/local/sbin/update-route53-arr

[Install]
WantedBy=multi-user.target

